# ops/launchd/com.joehu.market.premarket-options-intelligence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.joehu.market.premarket-options-intelligence</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/joehu/proj/market &amp;&amp; MARKET_DUCKDB_PATH=/Users/joehu/proj/market/data/investment.duckdb /opt/homebrew/bin/uv run python -m investment_panel.jobs.premarket_options_intelligence --config config.yaml</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/market-premarket-options-intelligence.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/market-premarket-options-intelligence.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/joehu/proj/market</string>
</dict>
</plist>
